# etc/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# ClamAV setting
0 10	* * *	root	/usr/bin/freshclam
0 20	* * 5	root	/usr/sbin/virusscan

#
